# roles/kube/files/zram-dirs
#!/sbin/openrc-run

depend() {
    need localmount
    before k0scontroller logger 
}

start() {
    ebegin "Mounting dirs to zram"
    /usr/local/bin/zram-dir /var/lib/k0s/etcd 512M
    /usr/local/bin/zram-dir /var/log 512M
    eend $?
}

stop() {
    ebegin "Unmounting dirs from zram"
    /usr/local/bin/zram-dir /var/lib/k0s/etcd
    /usr/local/bin/zram-dir /var/log
    eend $?
}
